# ucentral-client/rootfs/usr/lib/opkg/info/atfpolicy.list
/usr/sbin/atfpolicy
/etc/config/atfpolicy
/etc/init.d/atfpolicy
