# ucentral-client/rootfs/usr/lib/opkg/info/bridger.list
/lib/bpf/bridger-bpf.o
/etc/config/bridger
/etc/init.d/bridger
/usr/sbin/bridger
