# ucentral-client/rootfs/usr/lib/opkg/info/firstcontact.list
/usr/sbin/digicert
/etc/uci-defaults/zzz-firstcontact
/etc/init.d/firstcontact
/usr/share/ucentral/firstcontact.uc
