# ucentral-client/rootfs/usr/lib/opkg/info/ieee8021x.list
/usr/bin/ieee8021x.uc
/etc/config/ieee8021x
/etc/init.d/ieee8021x
