# ucentral-client/rootfs/usr/lib/opkg/info/radius-gw-proxy.list
/etc/init.d/radius-gw-proxy
/usr/sbin/radius-gw-proxy
