# ucentral-client/rootfs/usr/lib/opkg/info/radsecproxy.list
/etc/init.d/radsecproxy
/usr/sbin/radsecproxy
/etc/config/radsecproxy
