# ucentral-client/rootfs/usr/lib/opkg/info/ratelimit.list
/usr/bin/ratelimit
/etc/config/ratelimit
/etc/init.d/ratelimit
