# ucentral-client/rootfs/usr/lib/opkg/info/spotfilter.list
/usr/sbin/spotfilter
/etc/hotplug.d/net/10-spotfilter
/etc/init.d/spotfilter
/lib/bpf/spotfilter-bpf.o
