# ucentral-client/rootfs/usr/lib/opkg/info/switch-fabric.list
/usr/libexec/switch.sh
/etc/init.d/switch-fabric
/etc/config/switch
