# ucentral-client/rootfs/usr/lib/opkg/info/tip-defaults.list
/etc/ucentral/cas.pem
/etc/init.d/tip-version
