# ucentral-client/rootfs/usr/lib/opkg/info/ucentral-event.list
/etc/events.json
/etc/config/event
/usr/sbin/ucentral-wifiscan
/usr/sbin/ucentral-event
/etc/init.d/ucentral-event
/etc/init.d/ucentral-wifiscan
