# ucentral-client/rootfs/usr/lib/opkg/info/ucentral-dataplane.list
/usr/libexec/dataplane.sh
/etc/config/dataplane
/etc/hotplug.d/iface/80-dataplane
/etc/init.d/dataplane
