# ucentral-client/rootfs/usr/lib/opkg/info/ucentral-schema.list
/usr/share/ucentral/templates/services/gps.uc
/etc/ucentral/examples/ssh.json
/usr/share/ucentral/templates/services/quality_of_service.uc
/usr/share/ucentral/state.uc
/usr/share/ucentral/templates/metric/telemetry.uc
/etc/ucentral/examples/owe-transition.json
/etc/ucentral/examples/quality-threshold.json
/usr/share/ucentral/templates/services/igmp.uc
/usr/share/ucentral/capabilities.uc
/etc/ucentral/examples/dual-stack.json
/etc/ucentral/examples/vxlan.json
/usr/share/ucentral/cmd_leds.uc
/usr/share/ucentral/ip-collide.uc
/usr/share/ucentral/templates/interface/mesh.uc
/usr/libexec/ucentral/ucentral_led_blink.sh
/etc/ucentral/examples/fingerprint-raw.json
/etc/ucentral/examples/unit.json
/etc/sysctl.d/20-oom-reboot.conf
/etc/ucentral/examples/dhcp-relay.json
/usr/share/ucentral/templates/interface/broadband.uc
/etc/ucentral/examples/mesh.json
/etc/uci-defaults/99-ucentral-opennds
/usr/share/ucentral/diagnostic.uc
/usr/share/ucentral/wifi/mesh.uc
/etc/ucentral/examples/wwan.json
/usr/share/ucentral/cmd_syslog.uc
/etc/ucentral/examples/switch-fabric.json
/usr/share/ucentral/templates/interface/gre.uc
/usr/share/ucentral/templates/metric/statistics.uc
/usr/share/ucentral/templates/services/captive.uc
/etc/ucentral/examples/switch-vlan.json
/usr/share/ucentral/cmd_rtty.uc
/etc/ucentral/examples/onboarding-ap.yml
/usr/share/ucentral/reboot_cause.uc
/etc/ucentral/examples/crypto-enterprise.json
/usr/share/ucentral/templates/services/dhcp_snooping.uc
/usr/share/ucentral/cmd_stats.uc
/usr/share/ucentral/templates/interface/firewall/allow.uc
/etc/ucentral/examples/qos-class.json
/usr/share/ucentral/templates/toplevel.uc
/etc/uci-defaults/99-ucentral-event
/usr/share/ucentral/wifi/phy.uc
/etc/ucentral/examples/block-rfc1918.json
/etc/uci-defaults/99-ucentral-qos
/usr/share/ucentral/templates/services/ssh.uc
/etc/ucentral/examples/big.json
/etc/ucentral/examples/rate-limit.json
/etc/ucentral/examples/captive-click.json
/usr/share/ucentral/wifi/iface.uc
/usr/share/ucentral/templates/services/lldp.uc
/etc/ucentral/examples/ieee8021x-radius.json
/usr/share/ucentral/templates/interface/common.uc
/usr/share/ucentral/cmd_factory.uc
/usr/share/ucentral/templates/services/wireguard_overlay.uc
/etc/ucentral/examples/radius-proxy.json
/etc/uci-defaults/99-ucentral-umdns
/usr/share/ucentral/qos.json
/etc/ucentral/examples/wds-ap.json
/usr/share/ucentral/templates/ethernet.uc
/etc/uci-defaults/99-ucentral-rtty
/usr/share/ucentral/templates/services/ntp.uc
/usr/share/ucentral/templates/metric/wifi_scan.uc
/etc/ucentral/examples/captive-webroot.json
/usr/share/ucentral/wifi_max_user.uc
/usr/bin/tcpdump_timeout
/usr/share/ucentral/templates/interface/ieee8021x.uc
/usr/share/ucentral/templates/services/radius_proxy.uc
/etc/ucentral/examples/ieee8021x-nat.json
/etc/ucentral/examples/ntp.json
/etc/ucentral/ucentral.cfg.0000000001
/etc/ucentral/examples/metrics.json
/etc/ucentral/examples/radius-request-attr.json
/etc/ucentral/examples/rrm.json
/etc/ucentral/examples/ieee8021x-mac-auth.json
/etc/ucentral/examples/fingerprint-final-periodic.json
/usr/share/ucode/bundle.uc
/etc/ucentral/schema.json
/usr/share/ucentral/templates/services/data_plane.uc
/etc/ucentral/examples/owe.json
/etc/ucentral/examples/captive.json
/usr/libexec/naptr_lookup.sh
/etc/ucentral/examples/online-check.json
/usr/share/ucentral/templates/services/rrm.uc
/etc/ucentral/examples/crypto-psk.json
/etc/uci-defaults/99-ucentral-network
/usr/share/ucentral/templates/services/radius_gw_proxy.uc
/usr/share/ucentral/templates/interface/firewall/forward.uc
/usr/share/ucentral/templates/metric/dhcp_snooping.uc
/usr/share/ucentral/templates/unit.uc
/etc/ucentral/examples/wpa2-radius.json
/etc/ucentral/examples/vlan.json
/etc/ucentral/examples/switch-ports.json
/etc/ucentral/examples/wds-sta.json
/usr/share/ucentral/cmd.uc
/usr/share/ucentral/telemetry.uc
/etc/ucentral/examples/fingerprint-final-always.json
/etc/ucentral/examples/loop-detect.json
/usr/share/ucentral/templates/interface/gre6.uc
/usr/share/ucentral/templates/interface/firewall.uc
/usr/share/ucentral/cmd_transfer.uc
/etc/ucentral/examples/wifi-6e-afc.json
/usr/share/ucentral/templates/services/wifi_steering.uc
/usr/share/ucentral/templates/interface/ipv4.uc
/usr/share/ucentral/templates/services/rtty.uc
/usr/share/ucentral/templates/services/log.uc
/etc/ucentral/examples/captive-multiple.json
/etc/ucentral/examples/default.json
/etc/ucentral/examples/gps.json
/etc/ucentral/examples/qos.json
/usr/share/ucentral/templates/admin_ui.uc
/etc/ucentral/examples/psk2-radius.json
/etc/ucentral/examples/data-plane.json
/etc/ucentral/examples/radius-gw-proxy.json
/etc/ucentral/examples/admin_ui.json
/etc/ucentral/examples/captive-credentials.json
/etc/ucentral/examples/radius.json
/usr/share/ucentral/cmd_trace.uc
/usr/share/ucentral/templates/metric/realtime.uc
/etc/ucentral/examples/captive-uam.json
/usr/libexec/ucentral/ucentral_apply.sh
/etc/config/onlinecheck
/etc/ucentral/examples/multi-psk.json
/usr/share/ucentral/templates/interface/dhcp.uc
/usr/share/ucentral/schemareader.uc
/etc/ucentral/examples/radius-secondary.json
/usr/share/ucentral/cmd_perform.uc
/usr/share/ucentral/templates/services/airtime_fairness.uc
/etc/uci-defaults/99-ucentral-wifi
/etc/ucentral/examples/dhcpsnoop.json
/usr/share/ucentral/onlinecheck.uc
/usr/share/ucentral/templates/eap_users.uc
/etc/ucentral/examples/captive-radius.json
/usr/share/ucentral/templates/base.uc
/usr/share/ucentral/templates/broadband.uc
/etc/ucentral/examples/igmp.json
/etc/ucentral/examples/strict-forwarding.json
/usr/share/ucentral/signature.uc
/usr/share/ucentral/templates/interface/ipv6.uc
/usr/share/ucentral/cmd_wifiscan.uc
/etc/uci-defaults/99-ucentral-hostname
/usr/share/ucentral/templates/services/ieee8021x.uc
/etc/ucentral/examples/gre.json
/etc/ucentral/examples/tip-oui.json
/etc/uci-defaults/99-ucentral-lldp
/usr/share/ucentral/templates/services/mdns.uc
/etc/uci-defaults/99-ucentral-uhttpd
/usr/share/ucentral/sysinfo.uc
/usr/share/ucentral/templates/interface/l2tp.uc
/usr/share/ucentral/wifi/survey.uc
/usr/share/ucentral/templates/interface/wmm.uc
/usr/share/ucentral/templates/interface.uc
/etc/ucentral/examples/steering.json
/usr/share/ucentral/templates/interface/ssid.uc
/usr/share/ucentral/templates/services/http.uc
/usr/share/ucentral/templates/interface/vxlan.uc
/usr/share/ucentral/templates/metric/health.uc
/etc/ucentral/examples/telemetry.json
/etc/board.d/04-regdm
/etc/ucentral/examples/lldp.json
/usr/share/ucentral/templates/config_raw.uc
/usr/share/ucentral/templates/metric/wifi_frames.uc
/usr/share/ucentral/templates/services/dhcp_relay.uc
/usr/share/ucentral/health.uc
/usr/share/ucentral/templates/services/fingerprint.uc
/usr/share/ucentral/probe_services.uc
/etc/ucentral/examples/maverick.json
/usr/share/ucentral/templates/switch.uc
/usr/share/ucentral/renderer.uc
/etc/ucentral/examples/wifi-7.json
/usr/share/ucentral/templates/services/online_check.uc
/usr/share/ucentral/cmd_fingerprint.uc
/usr/share/ucentral/cmd_upgrade.uc
/usr/share/ucentral/cmd_certupdate.uc
/etc/init.d/onlinecheck
/usr/share/ucentral/templates/spotfilter.uc
/etc/ucentral/examples/eap_local.json
/usr/share/ucentral/ucentral.uc
/usr/share/ucentral/wifi/station.uc
/etc/ucentral/examples/gre6.json
/usr/share/ucentral/templates/interface/captive.uc
/usr/share/ucentral/templates/services/facebook_wifi.uc
/usr/share/ucentral/cmd_reboot.uc
/usr/share/ucentral/cmd_rrm.uc
/etc/ucentral/examples/wifi-6e.json
/usr/share/ucentral/cmd_script.uc
/usr/share/ucentral/templates/interface/bridge-vlan.uc
/etc/ucentral/examples/single-ssid.yml
/usr/share/ucentral/templates/radio.uc
