# ucentral-client/rootfs/usr/lib/opkg/info/ucentral-state.list
/usr/sbin/ucentral-pstore
/etc/config/state
/usr/sbin/ucentral-state
/etc/init.d/ucentral-state
