# ucentral-client/rootfs/usr/lib/opkg/info/ucentral-client.list
/etc/config/ucentral
/usr/sbin/ucentral
/etc/config/ustats
/etc/ucentral/version.json
/etc/uci-defaults/zzz-ucentral
/etc/init.d/ucentral
