# ucentral-client/rootfs/usr/lib/opkg/info/udevstats.list
/etc/config/udevstats
/etc/init.d/udevstats
/lib/bpf/udevstats.o
/etc/hotplug.d/net/10-udevstats
/usr/sbin/udevstats
