# ucentral-client/rootfs/usr/lib/opkg/info/udevmand.list
/usr/sbin/udevmand
/etc/init.d/udevmand
