# ucentral-client/rootfs/usr/lib/opkg/info/udhcpsnoop.list
/etc/hotplug.d/net/10-dhcpsnoop
/etc/config/dhcpsnoop
/etc/init.d/dhcpsnoop
/usr/sbin/udhcpsnoop
