# ucentral-client/rootfs/usr/lib/opkg/info/udhcprelay.list
/etc/init.d/dhcprelay
/usr/sbin/udhcprelay
/etc/hotplug.d/net/10-dhcprelay
/etc/config/dhcprelay
