# ucentral-client/rootfs/usr/lib/opkg/info/udnssnoop.list
/usr/sbin/udnssnoop
/etc/init.d/dnssnoop
