# ucentral-client/rootfs/usr/lib/opkg/info/ugps.list
/etc/config/gps
/etc/init.d/ugps
/usr/sbin/ugps
