# ucentral-client/rootfs/usr/lib/opkg/info/umdns.list
/etc/config/umdns
/usr/sbin/umdns
/etc/init.d/umdns
/etc/seccomp/umdns.json
