# ucentral-client/rootfs/usr/lib/opkg/info/unetd.list
/lib/netifd/proto/unetd.sh
/lib/bpf/mss.o
/usr/sbin/unetd
/etc/init.d/unetd
