# ucentral-client/rootfs/usr/lib/opkg/info/uspot.list
/www-uspot/images/OpenWiFi_BadgeLogo_White.svg
/usr/share/uspot/dump-env.uc
/usr/share/uspot/cpd.uc
/usr/share/uspot/header
/usr/share/uspot/firewall.nft
/usr/share/uspot/firewall.ipt
/usr/share/uspot/handler.uc
/usr/share/uspot/redir.uc
/usr/lib/ucode/uam.so
/usr/share/uspot/auth.uc
/usr/share/uspot/allow.uc
/usr/share/uspot/footer
/www-uspot/images/OpenWiFi_LogoLockup_DarkGreyColour.svg
/usr/share/uspot/portal.uc
/usr/share/uspot/uspot.uc
/www-uspot/css/captive-portal.css
/usr/share/uspot/connected.uc
/etc/config/uspot
/usr/share/uspot/handler-env.uc
/usr/share/uspot/handler-uam.uc
/usr/share/uspot/logoff.uc
/usr/share/uspot/credentials.uc
/usr/share/uspot/radius.uc
/usr/share/uspot/click.uc
/usr/bin/radius-client
/www-uspot/css/bootstrap.min.css
/usr/share/uspot/handler-cpd.uc
/usr/bin/captive
/etc/init.d/uspot
/usr/share/uspot/error.uc
