# ucentral-client/rootfs/usr/lib/opkg/info/usteer.list
/sbin/usteerd
/sbin/ap-monitor
/usr/libexec/uchannel.uc
/etc/init.d/usteer
/etc/config/usteer
