# ucentral-client/rootfs/usr/lib/opkg/info/ustp.list
/etc/init.d/ustpd
/sbin/ustpd
/sbin/bridge-stp
